# supervisord.conf
[program:php-worker]
directory=/home/vagrant/Worker-Test/
command=php worker.php
process_name=%(program_name)s_%(process_num)02d
numprocs=3
autostart=true
